# crontab
* * * * * /root/test.sh
1,31 * * * * /root/script.sh
58 23 * * * /root/reboot.sh
