# crontab
* * * * * /root/kill.sh
1,31 * * * * /root/crashdetect.sh
58 23 * * * /root/reboot.sh
